# doquerinhos_shell/var/lib/dpkg/info/initscripts.list
/.
/sbin
/sbin/fsck.nfs
/etc
/etc/init.d
/etc/init.d/checkroot.sh
/etc/init.d/checkroot-bootclean.sh
/etc/init.d/hostname.sh
/etc/init.d/umountroot
/etc/init.d/bootmisc.sh
/etc/init.d/mountnfs-bootclean.sh
/etc/init.d/mountnfs.sh
/etc/init.d/mountdevsubfs.sh
/etc/init.d/halt
/etc/init.d/mountkernfs.sh
/etc/init.d/ondemand
/etc/init.d/reboot
/etc/init.d/umountnfs.sh
/etc/init.d/umountfs
/etc/init.d/checkfs.sh
/etc/init.d/sendsigs
/etc/init.d/single
/etc/init.d/rc.local
/etc/init.d/mountall-bootclean.sh
/etc/init.d/killprocs
/etc/init.d/urandom
/etc/init.d/skeleton
/etc/init.d/mountall.sh
/etc/default
/etc/default/halt
/etc/default/rcS
/etc/default/devpts
/lib
/lib/init
/lib/init/vars.sh
/var
/var/lib
/var/lib/initscripts
/var/lib/urandom
/var/log
/var/log/fsck
/run
/sys
/usr
/usr/share
/usr/share/lintian
/usr/share/lintian/overrides
/usr/share/lintian/overrides/initscripts
/usr/share/doc
/usr/share/doc/initscripts
/usr/share/doc/initscripts/copyright
/usr/share/doc/initscripts/README.Debian
/usr/share/doc/initscripts/NEWS.Debian.gz
/usr/share/man
/usr/share/man/man8
/usr/share/man/man8/fsck.nfs.8.gz
/usr/share/man/man5
/usr/share/man/man5/halt.5.gz
/usr/share/man/man5/rcS.5.gz
/usr/share/doc/initscripts/changelog.Debian.gz
